# cookbooks/wd42-site/templates/default/wd42-site.conf.erb
# The following pulled from:
# http://kvz.io/blog/2009/12/15/run-nodejs-as-a-service-on-ubuntu-karmic/

# cat /etc/init/transloaditapi2.conf
# http://upstart.ubuntu.com/wiki/Stanzas

description "wd42-site node.js server"
author      "Josh Gillies"

start on started static-network-up

post-start script
  # Fun things can go here...
end script

stop on shutdown
respawn
respawn limit 20 5

# Max open files are @ 1024 by default. Bit few.
limit nofile 32768 32768

env HOME=<%= node[:wd42_site][:dir] %>

script
  export HOME=$HOME
  chdir $HOME

  set -e
  mkfifo /tmp/wd42-site-log-fifo
  ( logger -t wd42-site </tmp/wd42-site-log-fifo & )
  exec >/tmp/wd42-site-log-fifo
  rm /tmp/wd42-site-log-fifo
  exec sudo node /srv/www/wd42-site/app.js 2>&1 >> /var/log/node.log
end script
